# templates/etc/systemd/system/etc-resolv-conf-watcher.path.j2
[Path]
Unit=etc-resolv-conf-watcher.path.service
PathChanged=/etc/resolv.conf

[Install]
WantedBy=multi-user.target
