# templates/etc/systemd/system/etc-resolv-conf-watcher.path.j2
[Path]
Unit=etc-resolv-conf-watcher.service
PathChanged=/run/resolvconf/resolv.conf

[Install]
WantedBy=multi-user.target
